# containers/fake-rootfs/lib/systemd/system/packagekit-offline-update.service
[Unit]
Description=Updates the operating system whilst offline
Requires=dbus.socket
OnFailure=reboot.target

[Service]
ExecStart=/usr/lib/packagekit/pk-offline-update
